# src/disk1/etc/systemd/user/shell.service
[Unit]
Description=Get shell service for user